# supervisord.conf
[program:mibot]
command=/PATH_TO/xmpp.py
directory=/DEFAULT_PATH
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/tmp/bot.out.log
stderr_logfile=/tmp/bot.err.log
